# gn/services/archive-pubmed.timer
[Unit]
Description=Daily PubMed downloads

[Timer]
OnBootSec=5h
OnUnitInactiveSec=15h

[Install]
WantedBy=default.target
