# gn/services/gitea-dump.timer
[Unit]
Description=Backup Gitea git server timer

[Timer]
OnBootSec=15m
OnUnitInactiveSec=1day

[Install]
WantedBy=default.target
